# srv/salt/grafana/files/grafana.service.jinja
[Unit]
Description=Grafana Container (Podman)
After=container-prometheus.service

[Service]
Type=simple
ExecStartPre=-/usr/bin/podman rm -f grafana
ExecStart=/usr/bin/podman run \
    --name grafana \
    --rm \
    --network {{ network_name }} \
    --ip {{ grafana_ip }} \
    -v {{ grafana_conf }}/grafana.ini:/etc/grafana/grafana.ini:ro,Z \
    -v {{ grafana_conf }}/provisioning:/etc/grafana/provisioning:ro,Z \
    -v {{ grafana_data }}:/var/lib/grafana:Z \
    -e GF_SECURITY_ADMIN_PASSWORD={{ admin_password }} \
    --health-cmd="wget -qO /dev/null http://localhost:3000/api/health || exit 1" \
    --health-interval=30s \
    --health-timeout=5s \
    --health-retries=3 \
    --stop-timeout=15 \
    {{ image }}
ExecStop=/usr/bin/podman stop -t 15 grafana
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
