# srv/salt/wp_maintenance/files/wp-maintenance.timer.jinja
[Unit]
Description=WordPress Maintenance Timer

[Timer]
OnCalendar=*-*-* {{ hour }}:{{ minute }}:00
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
