# srv/salt/backup/files/wp-backup.service.jinja
[Unit]
Description=WordPress Full Backup

[Service]
Type=oneshot
ExecStart={{ wp_bin }}/backup.sh
