# srv/salt/backup/files/wp-backup.timer.jinja
[Unit]
Description=Daily WordPress Backup Timer

[Timer]
OnCalendar=*-*-* {{ '%02d' % schedule_hour }}:{{ '%02d' % schedule_minute }}:00
Persistent=true

[Install]
WantedBy=timers.target
